# sbackup/etc/cron.d/sbackup
# /etc/cron.d/mentes

0 12,15,20 * * *     root   /usr/local/bin/sbackup
0 13 * * 5,6        root   //usr/local/bin/sbackup-archiv

#
