# sbackup/etc/cron.d/sbackupnoexcl
# /etc/cron.d/mentes: crontab fragment for save important data

30 20 * * 6     root   /usr/local/bin/sbackup
